# luci-app-flowoffload/root/etc/uci-defaults/luci-flowoffload
uci batch <<EOF
	delete ucitrack.@flowoffload[-1]
	add ucitrack flowoffload
	set ucitrack.@flowoffload[-1].init=flowoffload
	commit ucitrack
EOF
